# mysql/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql-server]
command=mysqld_safe
